# partition/roles/sflow-collector/files/sflow-mmdb-refresh.timer
[Unit]
Description=Monthly refresh of DB-IP GeoIP and ASN mmdb files

[Timer]
OnCalendar=monthly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
