# partition/roles/sflow-collector/files/vrf-port-export.service
[Unit]
Description=VRF/VNI to interface mapping exporter for the sFlow collector (ifindex enrichment)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/vrf-port-export.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
