# partition/roles/monitoring/bgp-metrics/templates/bgp-metrics.service.j2
[Unit]
Description=Write BGP metrics to the node-exporter textfile collector

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bgp-metrics.py
