# partition/roles/sonic/files/frr-reload.service
[Unit]
Description=Save and reload ConfigDB

[Service]
Type=oneshot
ExecStart=/usr/bin/docker exec bgp /usr/lib/frr/frr-reload
StandardOutput=journal

[Install]
WantedBy=multi-user.target