# partition/roles/metal-core/files/bgp-neighbor-stats.timer
[Unit]
Description=Collect BGP Neighbor statistics

[Timer]
OnCalendar=*:0/2
Persistent=false

[Install]
WantedBy=timers.target
