# partition/roles/metal-core/files/bgp-neighbor-stats.service
[Unit]
Description=Collect BGP Neighbor statistics

[Service]
Type=simple
ExecStart=/usr/bin/bgp-neighbor-stats.sh